# ecsagent.service
[Unit]
Description=ecsagent
After=network.target

[Service]
ExecStart=/usr/local/bin/ecsagent

[Install]
WantedBy=multi-user.target